# zomboid-updater.timer
[Unit]
Description=Run Project Zomboid update check every 10 minutes
Requires=zomboid-updater.service

[Timer]
# Run 2 minutes after boot (give Docker time to settle)
OnBootSec=2min
# Then every 10 minutes
OnUnitActiveSec=10min
# If a run was missed (e.g. system was off), catch up immediately
Persistent=true

[Install]
WantedBy=timers.target
